# build/controlone-agent.service
[Unit]
Description=Control One Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/controlone-agent --config /etc/control-one/nodeagent.yaml
Restart=on-failure
RestartSec=5
User=root
Group=root
LimitNOFILE=65536
AmbientCapabilities=CAP_SYSLOG CAP_DAC_OVERRIDE
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/control-one/nodeagent /var/log/control-one/nodeagent

[Install]
WantedBy=multi-user.target
